# raspbian-addons/etc/network/if-pre-up.d/iptables
#!/bin/sh
#restore iptables on boot

iptables-restore < /etc/iptables.up.rules
